# l/jre/profile.d/jre.sh
#!/bin/sh
export JAVA_HOME=/usr/lib64/java
export MANPATH="$MANPATH:/usr/lib64/java/man"
export PATH="$PATH:/usr/lib64/java/bin:/usr/lib64/java/jre/bin"
